# iso/airootfs/etc/systemd/system/greetd.service
[Unit]
Description=KOOMPI OS Greeter (greetd)
After=systemd-user-sessions.service plymouth-quit-wait.service getty@tty1.service
Conflicts=getty@tty1.service

[Service]
Type=simple
ExecStart=/usr/bin/greetd
Restart=always
RestartSec=1
Environment=XDG_SESSION_TYPE=wayland

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
