# iso/airootfs/etc/systemd/system/koompi-shell.service
[Unit]
Description=KOOMPI Shell (Rust Compositor)
After=network.target systemd-user-sessions.service
Conflicts=getty@tty1.service

[Service]
ExecStart=/usr/local/bin/koompi-shell
Restart=always
User=root
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
